# deployment/systemd/feishu-bot.service
[Unit]
Description=Feishu Bot Service (Bot + Web Admin Interface)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/feishu-bot
Environment="PATH=/opt/feishu-bot/venv/bin"
EnvironmentFile=/opt/feishu-bot/.env

# Start the bot service (includes both bot and web admin)
ExecStart=/opt/feishu-bot/venv/bin/python lark_bot.py

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/feishu-bot/logs /opt/feishu-bot/data

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
